# crontab
* * * * * /home/app/crawler/check-for-updates.sh
